# meta-agenc/recipes-support/agenc-ota/files/agenc-ota-check.timer
[Unit]
Description=AgenC OTA Update Check Timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=6h
RandomizedDelaySec=30min

[Install]
WantedBy=timers.target
