# meta-agenc/recipes-core/systemd/agenc-runtime.service
[Unit]
Description=AgenC ONE Agent Runtime
After=network-online.target sound.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=agenc
Group=agenc
WorkingDirectory=/opt/agenc
ExecStart=/usr/bin/python3 -u /opt/agenc/agenc_voice_task.py
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=10
WatchdogSec=120

# Environment
Environment=HOME=/data/agenc
Environment=XAI_API_KEY=
Environment=AGENC_WALLET_PATH=/data/keystore/wallet.json
Environment=AGENC_TASKS_FILE=/data/agenc/tasks.json
Environment=AGENC_LOG_FILE=/data/agenc/agent.log
EnvironmentFile=-/data/agenc/env

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/data
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=no

# Resource limits
MemoryMax=300M
CPUQuota=90%

# Hardware access
SupplementaryGroups=gpio spi i2c audio video

[Install]
WantedBy=multi-user.target
