# meta-agenc/recipes-security/dm-crypt/files/keystore-unlock.service
[Unit]
Description=Unlock AgenC Encrypted Keystore
After=local-fs.target
Before=agenc-runtime.service
ConditionPathExists=/dev/mmcblk0p5

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh -c 'SERIAL=$(cat /proc/cpuinfo | grep Serial | awk "{print \\$3}") && echo -n "agenc-$SERIAL" | sha256sum | awk "{print \\$1}" | tr -d "\\n" | cryptsetup luksOpen /dev/mmcblk0p5 agenc-keystore -'
ExecStart=/bin/mount /dev/mapper/agenc-keystore /data/keystore
ExecStop=/bin/umount /data/keystore
ExecStop=/sbin/cryptsetup luksClose agenc-keystore

[Install]
WantedBy=multi-user.target
